# docs/Ops/crontab/api
###app10-115
<pre>
 #把昨天的api日志同步到hadoop
 01 0 * * * /bin/bash /data1/mobile-log/cp_apilog_hadoop.sh 2>&1 >> /data1/mobile-log/cp_to_hadoop.log
 
 05 * * * * /bin/bash /data1/mobile-log/cp_apilog_hadoop_hour.sh 2>&1 >> /data1/mobile-log/cp_to_hadoop_hour.log

 #*/1 * * * * /bin/bash /home/www/publish-shell/jobs/activity_push_chat.sh
 
 30 1 * * * /bin/bash /home/www/bin/action_log_arch.sh >> /data1/mobile-log/action_log_arch.log
</pre>

###app10-116
<pre>
 #* * * * * /usr/local/redis/bin/redis-cli -h 10.10.3.120 -p 1702 LTRIM chat_register_queue 0 0 > /dev/null
 
 #* * * * * /usr/local/redis/bin/redis-cli -h 10.10.3.120 -p 1702 LTRIM chat_register_device_queue 0 0 > /dev/null
 
 #* * * * * /usr/local/redis/bin/redis-cli -h 10.10.3.120 -p 1702 LTRIM queue_device_info_persistance 0 0 > /dev/null
  
 #0 2 * * * /home/www/bin/log-utils/put_push_log_to_ods.sh > /dev/null
  
 * * * * * /home/www/bin/monitor/push_knowing_jvm.sh > /dev/null

 
 10 2 * * * /home/www/bin/log-utils/clean_push.sh > /dev/null
 
 * * * * * /data1/mobileops/mobile-bin/monitor/push-netstat.sh > /dev/null
 
 #* * * * * /bin/bash /home/www/bin/monitor/amtp.sh &> /dev/null
</pre>

###app10-117
<pre>
 05 * * * * /usr/bin/pkill -f Mobile_Job_DealServiceRecommend
 
 */1 * * * * /bin/bash /home/www/bin/check_job_162.sh
 
 * * * * * /usr/local/php-fpm/bin/php /home/www/config/v2/anjuke-chat/chat-job/launcher.php --class=Chat_Job_Monitor_APNS >> /dev/null
 
 * * * * * /bin/bash /home/www/bin/job_monitor.sh
 
 #* * * * * /bin/bash /home/www/bin/monitor/push_knowing_all.sh

 * * * * * /usr/local/php-fpm/bin/php /home/www/config/v2/anjuke-chat/chat-job/launcher.php --class=Chat_Job_PostKnowing >> /data1/logs/push/knowing.log 2>&1

 * * * * * /bin/bash /home/www/bin/monitor/redis2.sh
 
 * * * * * /usr/local/php-fpm/bin/php /home/www/config/v2/anjuke-chat/chat-job/launcher_beta.php --class=Chat_Job_Monitor_Chat --action=doGetMessageCountFromWeb >> /data1/logs/push/doGetMessageCountFromWeb.log
 
 * * * * * /usr/local/php-fpm/bin/php /home/www/config/v2/anjuke-chat/chat-job/launcher_beta.php --class=Chat_Job_Monitor_Chat --action=doGetRegisterUserCountFromWeb >> /data1/logs/push/doGetRegisterUserCountFromWeb.log

 * * * * * /usr/local/php-fpm/bin/php /home/www/config/v2/anjuke-chat/chat-job/launcher_beta.php --class=Chat_Job_Monitor_Chat --action=doGetRubbishMessageCount >> /data1/logs/push/doGetRubbishMessageCount.log

 * * * * * /usr/local/php-fpm/bin/php /home/www/config/v2/user-site/app-mobile-job/launcher.php --class=Mobile_Job_BrokerNewPropNotice --is-beta=1 >> /data1/logs/v2/Mobile_Job_BrokerNewPropNotice.log

 00 02 * * * /usr/local/php-fpm/bin/php /home/www/config/v2/user-site/app-mobile-job/launcher.php --class=Mobile_Job_MetroPropStat >> /data1/logs/v2/user-job/haozu/MetroPropStat.log

 00 03 * * * /usr/local/php-fpm/bin/php /home/www/config/v2/anjuke-chat/chat-job/launcher.php --class=Chat_Job_ChatBasicFix >> /data1/logs/push/Chat_Job_ChatBasicFix.log

 00 03 * * * /usr/local/php-fpm/bin/php /home/www/config/v2/anjuke-chat/chat-job/launcher.php --class=Chat_Job_TwoCodeFix >> /data1/logs/push/two_code_fix.log

 00 08 * * * /usr/local/php-fpm/bin/php /home/www/config/v2/anjuke-chat/chat-job/launcher.php --class=Chat_Job_Remind_SignBrokerReply >> /data1/logs/push/msg.log

 00 09 * * * /usr/local/php-fpm/bin/php /home/www/config/v2/anjuke-chat/chat-job/launcher.php --class=Chat_Job_Remind_RemindBrokerReply >> /data1/logs/push/msg.log

 00 10 * * * /usr/local/php-fpm/bin/php /home/www/config/v2/anjuke-chat/chat-job/launcher.php --class=Chat_Job_Remind_RemindBrokerByApns >> /data1/logs/push/apns.log

 00 11 * * * /usr/local/php-fpm/bin/php /home/www/config/v2/anjuke-chat/chat-job/launcher_beta.php --class=Chat_Job_Count_NoReplyBrokers --time=am >> /data1/logs/push/no_reply_brokers.log

 00 16 * * * /usr/local/php-fpm/bin/php /home/www/config/v2/anjuke-chat/chat-job/launcher_beta.php --class=Chat_Job_Count_NoReplyBrokers --time=pm >> /data1/logs/push/no_reply_brokers.log

 00 21 * * * /usr/local/php-fpm/bin/php /home/www/config/v2/anjuke-chat/chat-job/launcher.php --class=Chat_Job_Remind_NoticeBroker

 #*/1 * * * * /bin/bash /home/www/bin/weiliao_check_queue.sh >> /data1/logs/push/weiliao_queue.log

 */1 * * * * /bin/bash /home/www/bin/weiliao_check.sh >> /data1/logs/push/weiliao.log

 */1 * * * * /bin/bash /home/www/bin/clean_disk.sh

 */1 * * * * /usr/local/php-fpm/bin/php /home/www/config/v2/anjuke-chat/chat-job/launcher.php --class=Chat_Job_KnowingRedisMemory

 #community_db数据更新脚本，迁移自x###app10-022
 #17 7 * * * /usr/local/anjuke-php-fpm/bin/php /home/www/community-import-data/run.php all >>/home/www/community-import-data/run.log 2>&1

 00 06 * * * /usr/local/anjuke-php-fpm/bin/php /home/www/community-import-data/run.php all >>/home/www/community-import-data/run.log 2>&1

 #微店连续分享天数，推送微聊消息给经纪人
 #00 10 * * * /usr/local/anjuke-php-fpm/bin/php /home/www/config/v2/user-site/app-mobile-job/launcher.php --class=Mobile_Job_DealWeshopSharePush 2>&1

 #市场活动15012，经纪人推荐用户送话费，充值job，每天00点，10点跑一次
 #00 0,10 * * * /usr/local/anjuke-php-fpm/bin/php /home/www/config/v2/user-site/app-mobile-job/launcher.php --class=Mobile_Job_RechargeActivity0501 2>&1

 #市场推送 推送频次：每天11:30  每周五15:00
 #30 11 * * * /usr/local/anjuke-php-fpm/bin/php /home/www/config/v2/user-site/app-mobile-job/launcher.php --class=Mobile_Job_OldUserPushByDevice    >> /data1/logs/v2/user-job/pushchat/Mobile_Job_OldUserPushByDevice_20150525.log

 #0 15 * * 5 /usr/local/anjuke-php-fpm/bin/php /home/www/config/v2/user-site/app-mobile-job/launcher.php --class=Mobile_Job_OldUserPushByDevice    >> /data1/logs/v2/user-job/pushchat/Mobile_Job_OldUserPushByDevice_20150515.log

 #常驻服务
 #0 0 1 1 * /usr/local/anjuke-php-fpm/bin/php /home/www/config/v2/user-site/app-mobile-job/launcher.php --class=Mobile_Job_UserDeviceInfoPersistence --sleep-time=1

 #45 16 * * * /usr/local/anjuke-php-fpm/bin/php /home/www/config/v2/user-site/app-mobile-job/launcher.php --class=Mobile_Job_UserDeviceInfoPersistence --sleep-time=1

 #57 17 * * * /usr/local/anjuke-php-fpm/bin/php /home/www/config/v2/user-site/app-mobile-job/launcher.php --class=Mobile_Job_DealServicePush --sleep-time=1
 
 以下属于###app10-117 不过不在crontab中
      10 /usr/local/anjuke-php-fpm/bin/php /home/www/config/v2/anjuke-chat/chat-job/launcher.php --class=Chat_Job_Message_NotifyProcess --queue_name=chat_register_amtp_user_queue
      10 /usr/local/anjuke-php-fpm/bin/php /home/www/config/v2/anjuke-chat/chat-job/launcher.php --class=Chat_Job_Message_NotifyProcess --queue_name=chat_queue_app_notify
      10 /usr/local/anjuke-php-fpm/bin/php /home/www/config/v2/anjuke-chat/chat-job/launcher.php --class=Chat_Job_Message_NotifyProcess --queue_name=chat_queue_app_amtp_notify
       5 /usr/local/anjuke-php-fpm/bin/php /home/www/config/v2/user-site/app-mobile-job/launcher.php --class=Mobile_Job_JavaUserDeviceInfoPersistence --sleep-time=1
       5 /usr/local/anjuke-php-fpm/bin/php /home/www/config/v2/anjuke-chat/chat-job/launcher.php --class=Chat_Job_Message_NotifyProcess --queue_name=chat_queue_web_notify
       5 /usr/local/anjuke-php-fpm/bin/php /home/www/config/v2/anjuke-chat/chat-job/launcher.php --class=Chat_Job_Message_NotifyProcess --queue_name=chat_json_register_device_queue
       4 /usr/local/anjuke-php-fpm/bin/php /home/www/config/v2/user-site/app-mobile-job/launcher.php --class=Mobile_Job_DeviceInfoPersistence --sleep-time=1
       4 /usr/local/anjuke-php-fpm/bin/php /home/www/config/v2/user-site/app-mobile-job/launcher.php --class=Mobile_Job_DealServiceRecommend --sleep-time=1
       4 /usr/local/anjuke-php-fpm/bin/php /home/www/config/v2/anjuke-chat/chat-job/launcher.php --class=Chat_Job_Message_NotifyProcess --queue_name=chat_queue_account_message_batch_send
       3 /usr/local/anjuke-php-fpm/bin/php /home/www/config/v2/user-site/app-mobile-job/launcher.php --class=Mobile_Job_UserDeviceInfoPersistence --sleep-time=1
       2 /usr/local/anjuke-php-fpm/bin/php /home/www/config/v2/user-site/app-mobile-job/launcher.php --class=Mobile_Job_DealQaNewAnswer --sleep-time=1
       2 /usr/local/anjuke-php-fpm/bin/php /home/www/config/v2/user-site/app-mobile-job/launcher.php --class=Mobile_Job_B8 --sleep-time=1
       2 /usr/local/anjuke-php-fpm/bin/php /home/www/config/v2/anjuke-chat/chat-job/launcher.php --class=Chat_Job_Message_NotifyProcess --queue_name=chat_register_user_queue
       2 /usr/local/anjuke-php-fpm/bin/php /home/www/config/v2/anjuke-chat/chat-job/launcher.php --class=Chat_Job_Message_NotifyProcess --queue_name=chat_queue_device_login_notify
       2 /usr/local/anjuke-php-fpm/bin/php /home/www/config/v2/anjuke-chat/chat-job/launcher.php --class=Chat_Job_Message_NotifyProcess --queue_name=chat_queue_device_info
       2 /usr/local/anjuke-php-fpm/bin/php /home/www/config/v2/anjuke-chat/chat-job/launcher.php --class=Chat_Job_Message_NotifyProcess --queue_name=account_online_status_notify_http_chunk
       1 /usr/local/anjuke-php-fpm/bin/php /home/www/config/v2/anjuke-chat/chat-job/launcher.php --class=Chat_Job_Message_OfflinePushForIos
       1 /usr/local/anjuke-php-fpm/bin/php /home/www/config/v2/anjuke-chat/chat-job/launcher.php --class=Chat_Job_Message_OfflineDevicePushForIos
       1 /usr/local/anjuke-php-fpm/bin/php /home/www/config/v2/anjuke-chat/chat-job/launcher.php --class=Chat_Job_Message_NotifyProcess --queue_name=chat_queue_receive_message_amtp_reply
       1 /usr/local/anjuke-php-fpm/bin/php /home/www/config/v2/anjuke-chat/chat-job/launcher.php --class=Chat_Job_Message_NotifyProcess --queue_name=chat_queue_read_message_sync_notify
       1 /usr/local/anjuke-php-fpm/bin/php /home/www/config/v2/anjuke-chat/chat-job/launcher.php --class=Chat_Job_Message_NotifyProcess --queue_name=chat_queue_read_message_notify
       1 /usr/local/anjuke-php-fpm/bin/php /home/www/config/v2/anjuke-chat/chat-job/launcher.php --class=Chat_Job_Message_NotifyProcess --queue_name=chat_queue_read_broadcast_message_notify
       1 /usr/local/anjuke-php-fpm/bin/php /home/www/config/v2/anjuke-chat/chat-job/launcher.php --class=Chat_Job_Message_NotifyProcess --queue_name=chat_queue_read_app_message_notify
       1 /usr/local/anjuke-php-fpm/bin/php /home/www/config/v2/anjuke-chat/chat-job/launcher.php --class=Chat_Job_Message_NotifyProcess --queue_name=chat_queue_read_app_broadcast_message_notify
       1 /usr/local/anjuke-php-fpm/bin/php /home/www/config/v2/anjuke-chat/chat-job/launcher.php --class=Chat_Job_Message_NotifyProcess --queue_name=chat_queue_push_message_retry_for_ios
       1 /usr/local/anjuke-php-fpm/bin/php /home/www/config/v2/anjuke-chat/chat-job/launcher.php --class=Chat_Job_Message_NotifyProcess --queue_name=chat_queue_push_message_for_ios
       1 /usr/local/anjuke-php-fpm/bin/php /home/www/config/v2/anjuke-chat/chat-job/launcher.php --class=Chat_Job_Message_NotifyProcess --queue_name=chat_queue_push_device_message_for_ios
       1 /usr/local/anjuke-php-fpm/bin/php /home/www/config/v2/anjuke-chat/chat-job/launcher.php --class=Chat_Job_Message_NotifyProcess --queue_name=chat_queue_no_reply_sms_notify
       1 /usr/local/anjuke-php-fpm/bin/php /home/www/config/v2/anjuke-chat/chat-job/launcher.php --class=Chat_Job_Message_NotifyProcess --queue_name=chat_queue_new_user_message_process
       1 /usr/local/anjuke-php-fpm/bin/php /home/www/config/v2/anjuke-chat/chat-job/launcher.php --class=Chat_Job_Message_NotifyProcess --queue_name=chat_queue_new_service_message_process
       1 /usr/local/anjuke-php-fpm/bin/php /home/www/config/v2/anjuke-chat/chat-job/launcher.php --class=Chat_Job_Message_NotifyProcess --queue_name=chat_queue_new_device_message_queue
       1 /usr/local/anjuke-php-fpm/bin/php /home/www/config/v2/anjuke-chat/chat-job/launcher.php --class=Chat_Job_Message_NotifyProcess --queue_name=chat_queue_chat_basic_process
       1 /usr/local/anjuke-php-fpm/bin/php /home/www/config/v2/anjuke-chat/chat-job/launcher.php --class=Chat_Job_Message_NotifyProcess --queue_name=chat_queue_amtp_reply
       1 /usr/local/anjuke-php-fpm/bin/php /home/www/config/v2/anjuke-chat/chat-job/launcher.php --class=Chat_Job_Message_NotifyProcess --queue_name=account_online_status_notify_apns
       1 /usr/local/anjuke-php-fpm/bin/php /home/www/config/v2/anjuke-chat/chat-job/launcher.php --class=Chat_Job_Message_MonitorAckInfo
</pre>

###app10-125
<pre>
 30 1 * * * /bin/bash /home/www/bin/action_log_arch.sh >> /data1/mobile-log/action_log_arch.log
 
 35 1 * * * /bin/bash /home/www/publish-shell/gzip.sh >> /data1/mobile-log/gzip.log
</pre>

###app10-126
<pre>
 #* * * * * /home/www/bin/monitor/push_knowing_jvm.sh > /dev/null

 #* * * * * /usr/local/redis/bin/redis-cli -h 10.10.3.120 -p 1702 LTRIM chat_register_queue 0 0 > /dev/null

 #* * * * * /usr/local/redis/bin/redis-cli -h 10.10.3.120 -p 1702 LTRIM chat_register_device_queue 0 0 > /dev/null

 #* * * * * /usr/local/redis/bin/redis-cli -h 10.10.3.120 -p 1702 LTRIM queue_device_info_persistance 0 0 > /dev/null

 #0 2 * * * /home/www/bin/log-utils/put_push_log_to_ods.sh > /dev/null

 #10 2 * * * /home/www/bin/log-utils/clean_push.sh > /dev/null

 30 1 * * * /bin/bash /home/www/bin/action_log_arch.sh >> /data1/mobile-log/action_log_arch.log

 * * * * * /bin/bash /home/www/bin/monitor/amtp.sh > /dev/null

 0 1 * * * /home/www/bin/log-utils/put_amtp_log_to_ods.sh > /dev/null

 #* * * * * /usr/bin/rsync -ave ssh 10.10.3.91:/home/www/v2/ /home/www/v2/ >> /data1/logs/rsync.log 2>&1
</pre>

###app10-127
<pre>
 30 1 * * * /bin/bash /home/www/bin/action_log_arch.sh >> /data1/mobile-log/action_log_arch.log
 
 #*/1 * * * * /bin/bash /home/www/rsync-apk.sh

 01 06 * * * bash /data1/logs/slow_url/slow_url.sh >> /data1/logs/slow_url/mail.log 2>&1

 08 01 28 * * bash /home/www/publish-shell/delete_old_log.sh >> /data1/logs/clear_log.log 2>&1 &
</pre>
###app10-128
<pre>
 30 1 * * * /bin/bash /home/www/bin/action_log_arch.sh >> /data1/mobile-log/action_log_arch.log
 
 */1 * * * * /bin/bash /home/www/bin/weiliao_check_queue_new.sh >> /data1/logs/push/weiliao_queue_new.log
 
 #* * * * * /usr/local/anjuke-php-fpm/bin/php /home/www/config/v2/anjuke-chat/chat-job/launcher.php --class=Chat_Job_Monitor_APNS >>/tmp/x.log
</pre>
###app10-129
<pre>
 #add by bossni 敏感词服务词库更新检查重启脚本
 0 * * * * /bin/sh /home/www/bin/stringmatch_dict_sync.sh >>/data1/logs/v2/string_match_sync.log

 30 1 * * * /bin/bash /home/www/bin/action_log_arch.sh >> /data1/mobile-log/action_log_arch.log

 */1 * * * * /bin/bash /home/www/bin/weiliao_check_queue_new.sh >> /data1/logs/push/weiliao_queue_new.log
 </pre>
